# discord-bot/threat-scanner.service
[Unit]
Description=Conflict Globe Threat Scanner
After=network.target

[Service]
Type=simple
User=c0smic
WorkingDirectory=/home/c0smic/discord-bot
ExecStart=/usr/bin/node threat-scanner.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
