# antenna-agent/scripts/site-content.service
[Unit]
Description=Antenna Site Content Manager

[Service]
Type=oneshot
User=c0smic
WorkingDirectory=/home/c0smic/opencode/ui/antenna
ExecStart=/bin/bash /home/c0smic/opencode/ui/antenna/scripts/site-content.sh
StandardOutput=append:/home/c0smic/.antenna/logs/content-manager.log
StandardError=append:/home/c0smic/.antenna/logs/content-manager.log
