# antenna-agent/scripts/web-guardian.service
[Unit]
Description=Antenna Web Guardian - Monitors website and services

[Service]
Type=simple
User=c0smic
WorkingDirectory=/home/c0smic/opencode/ui/antenna
ExecStart=/bin/bash /home/c0smic/opencode/ui/antenna/scripts/web-guardian.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/c0smic/.antenna/logs/web-guardian.log
StandardError=append:/home/c0smic/.antenna/logs/web-guardian.log

[Install]
WantedBy=default.target
